# assets/systemd/kwybars-daemon.service
[Unit]
Description=Kwybars audio visualizer daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
Environment=PATH=%h/.cargo/bin:/usr/local/bin:/usr/bin
ExecStart=/usr/bin/env kwybars-daemon
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
